# var/lib/dpkg/info/php7.3-fpm.conffiles
/etc/apache2/conf-available/php7.3-fpm.conf
/etc/init.d/php7.3-fpm
/etc/init/php7.3-fpm.conf
/etc/logrotate.d/php7.3-fpm
/etc/php/7.3/fpm/php-fpm.conf
/etc/php/7.3/fpm/pool.d/www.conf
